# target/linux/fonera2-2.6/base-files/default/etc/uci-defaults/network
#!/bin/sh
if [ -e "/sys/bus/mdio_bus/drivers/Infineon ADM6996/0:00" -o \
     -e "/sys/bus/mdio_bus/drivers/Marvell 88E6060/0:1f" ]; then
	uci set network.lan.ifname=eth0.0
	uci set network.wan=interface
	uci set network.wan.ifname=eth0.1
	uci set network.wan.proto=dhcp
	uci commit network
fi
