# package/hotplug2/files/hotplug2.init
#!/bin/sh /etc/rc.common
START=65

boot() {
	#hotplug2 --persistent --coldplug &
	
	/sbin/usb-function-enable.sh &
}
